# roles/cfg_openwrt/files/corerouter/crontabs/root
18 * * * *	test -e /usr/sbin/owm.lua && /usr/sbin/owm.lua
*/10 * * * *    /etc/config/airmax_dfs_reset
