# zwingli/zwingli-core/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

/usr/sbin/ip -6 rule add from all lookup 192 prio 100

exit 0
